# src/systemd/qs-kwin-bridge.service
[Unit]
Description=Quickshell KDE Window Bridge
After=plasma-kwin_wayland.service

[Service]
ExecStart=/usr/bin/env python3 %h/.local/bin/qs-kwin-bridge.py
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
